# services.d/node/run
#!/command/with-contenv sh

node /usr/src/app/bin/www
